# packaging/deb/ubuntu/12.04/geowave/debian/geowave-jetty.postinst
#!/bin/bash

chown -fR 'geowave:geowave' '/usr/local/geowave/geoserver'

/usr/sbin/update-rc.d geowave defaults
